# sinapum-django.service
[Unit]
Description=MCP SinapUm - Master Control Program
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/MCP_SinapUm
Environment="PATH=/root/MCP_SinapUm/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/root/MCP_SinapUm"
Environment="DJANGO_SETTINGS_MODULE=setup.settings"
ExecStart=/root/MCP_SinapUm/.venv/bin/python3 -m gunicorn setup.wsgi:application --bind 0.0.0.0:5000 --workers 3 --timeout 120 --access-logfile - --error-logfile -
Restart=always
RestartSec=10

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sinapum-django

[Install]
WantedBy=multi-user.target

